# autostart/polybar.desktop
[Desktop Entry]
Type=Application
Exec=/home/alex/.config/polybar/launch.sh
Hidden=false
X-MATE-Autostart-enabled=true
Name=Polybar
